# systemd/adhan-daemon.timer
[Unit]
Description=Restart adhan-daemon daily to recalc times

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target